# deploy/course-xva-essentials.service
[Unit]
Description=Course XVA Essentials FastAPI backend
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/course_xva_essentials/backend
EnvironmentFile=-/var/www/html/course_xva_essentials/backend/.env
ExecStart=/var/www/html/course_xva_essentials/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8003 --proxy-headers --forwarded-allow-ips='*'
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
